# supervisord.conf
; supervisord.conf
; Manages two processes inside the Docker container:
;   1. nginx  — serves the React app on port 80
;   2. node   — runs the WebSocket server on port 3001
;
; supervisord acts as PID 1 so Docker can track the container's health.
; If either process crashes, supervisord restarts it automatically.
; This file is COPY'd into the Docker image at build time.

[supervisord]
; nodaemon=true keeps supervisord in the foreground so Docker doesn't exit
nodaemon=true
; Log to stdout so output appears in `docker logs`
logfile=/dev/null
pidfile=/tmp/supervisord.pid

[program:nginx]
; Run nginx in the foreground (daemon off) so supervisord can manage it
command=nginx -g "daemon off;"
autostart=true
autorestart=true
; Pipe nginx logs to Docker's stdout/stderr
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:node]
; Start the WebSocket server
command=node /app/server/server.js
autostart=true
autorestart=true
; Run as non-root for security
user=appuser
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0